# gw/client/etc/rc.local
bash /etc/setup-vpn.sh



bash /etc/openvpn/startup.sh


exit 0
